# agent/com.stwrd.agent-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd schedule for the STWRD background agent.

  Runs one loop per hour (only while the Mac is awake and online). Each run
  works the next queue item and leaves a claude/auto-<date> review branch; it
  never pushes or merges. When the queue is empty the run is a cheap no-op.

  Install:
    mkdir -p ~/.stwrd-agent
    cp agent/com.stwrd.agent-loop.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.stwrd.agent-loop.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.stwrd.agent-loop.plist
  Run once now (without waiting for the hour):
    launchctl start com.stwrd.agent-loop
  Watch:
    tail -f ~/.stwrd-agent/launchd.out.log
    git branch --list 'claude/auto-*'
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stwrd.agent-loop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vijayramdon/STWRD/scripts/agent-loop.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vijayramdon/STWRD</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/vijayramdon/.stwrd-agent/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vijayramdon/.stwrd-agent/launchd.err.log</string>
</dict>
</plist>
